# nikon.rules
SUBSYSTEMS=="usb", ATTRS{product}=="NIKON DSC D600", RUN+="/opt/photo/start-gphoto2"
